# supervisord.conf
[supervisord]
nodaemon=true
[program:sshd]
command=/usr/sbin/sshd -D
[program:squid]
command=/usr/sbin/squid -D
[program:crond]
command=crond -i
[program:iptables]
command=/usr/sbin/iptables -D
